# chien-treats/deployment/systemd/chienstreats-frontend.service
[Unit]
Description=Chien's Treats Frontend (Next.js)
After=network.target chienstreats-api.service
Wants=chienstreats-api.service

[Service]
Type=simple
User=chienstreats
Group=chienstreats
WorkingDirectory=/srv/sites/chienstreats/current
Environment=NODE_ENV=production
Environment=PORT=3101
EnvironmentFile=/srv/sites/chienstreats/shared/.env
ExecStart=/usr/bin/pnpm start
Restart=always
RestartSec=10s
StandardOutput=append:/srv/sites/chienstreats/shared/logs/frontend.log
StandardError=append:/srv/sites/chienstreats/shared/logs/frontend-error.log

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/srv/sites/chienstreats/shared/logs /srv/sites/chienstreats/shared/uploads
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true

[Install]
WantedBy=multi-user.target
