# chien-treats/deployment/systemd/chienstreats-api.service
[Unit]
Description=Chien's Treats API (NestJS)
After=network.target postgresql.service redis-server.service
Wants=postgresql.service redis-server.service

[Service]
Type=simple
User=chienstreats
Group=chienstreats
WorkingDirectory=/srv/sites/chienstreats/current
Environment=NODE_ENV=production
EnvironmentFile=/srv/sites/chienstreats/shared/.env.api
ExecStart=/usr/bin/node apps/api/dist/apps/api/src/main.js
Restart=always
RestartSec=10s
StandardOutput=append:/srv/sites/chienstreats/shared/logs/api.log
StandardError=append:/srv/sites/chienstreats/shared/logs/api-error.log

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/srv/sites/chienstreats/shared/logs /srv/sites/chienstreats/shared/uploads
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true

[Install]
WantedBy=multi-user.target
